# overlay-rootfs/etc/services.d/nebula/run
#!/bin/execlineb -P

backtick -E CONFIG_DIR { printcontenv CONFIG_DIR }
backtick -E _CONFIG_FILE {
	pipeline -d { find ${CONFIG_DIR} -name "*.y*ml" -exec stat -c "%Y %n" {} + }
	pipeline -d { s6-sort }
	awk "END{print $NF}"
}
backtick -D ${_CONFIG_FILE} -E CONFIG_FILE { printcontenv CONFIG_FILE }
ifelse -n {
	s6-setuidgid nebula
	redirfd -w 1 /dev/null
	nebula -test -config ${CONFIG_FILE}
} {
	foreground {
		fdmove -c 1 2
		s6-echo "[ERROR] Missing/invalid nebula config file. Aborting..."
	}
	getcwd -E PWD
	s6-svscanctl -tb "${PWD}/.."
}
s6-notifyoncheck -d
s6-setuidgid nebula
nebula -config ${CONFIG_FILE}
